# gentoo-purple-bridge.service
[Unit]
Description=Purple Bridge — container-to-host event relay
BindsTo=gentoo-nspawn.service
After=gentoo-nspawn.service
PartOf=gentoo-nspawn.service

[Service]
Type=simple
ExecStart=/usr/local/bin/purple-host-bridge.sh
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
