# systemd/system/pleiades-container.service
[Unit]
Description=Pleiades Gentoo systemd-nspawn substrate
Documentation=https://github.com/Zheke32174/pleiades-container
After=network-online.target local-fs.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=notify
NotifyAccess=all
Environment=PLEIADES_ROOT=/var/lib/machines/pleiades
EnvironmentFile=-/etc/pleiades/container.env
ExecStart=/usr/bin/systemd-nspawn --settings=no --keep-unit --boot --machine=pleiades --directory=${PLEIADES_ROOT} --network-veth --resolv-conf=copy-host
ExecStop=/usr/bin/machinectl poweroff pleiades
Restart=on-failure
RestartSec=15s
TimeoutStartSec=120s
TimeoutStopSec=45s
KillMode=mixed
Delegate=yes
TasksMax=1024
CPUQuota=200%
MemoryHigh=2G
MemoryMax=3G
OOMPolicy=stop

[Install]
WantedBy=multi-user.target
